# configs/services.d/openxpki/run
#!/usr/bin/execlineb -P
openxpkictl --foreground start