# configs/services.d/apache2/run
#!/usr/bin/execlineb -P
apache2ctl -DFOREGROUND